# playbooks/configs/etc/cron.daily/certbot.j2
#!/bin/bash
/usr/bin/certbot renew --quiet
/usr/sbin/service nginx reload